# deploy/systemd/iris-server.service
[Unit]
Description=Iris Server - 监控数据中心服务器
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=iris
Group=iris
WorkingDirectory=/opt/iris
Environment="RUST_LOG=info"
ExecStart=/opt/iris/iris-server --addr 0.0.0.0:50051
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/iris

# 资源限制
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
